# root/etc/s6-overlay/s6-rc.d/init-novnc-config/run
#!/usr/bin/with-contenv bash

# default file copies first run
if [[ ! -f /usr/share/novnc/index.html ]]; then
  ln -s /usr/share/novnc/vnc.html /usr/share/novnc/index.html
fi
